# supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock

[supervisord]
logfile  = /tmp/supervisord.log
pidfile  = /tmp/supervisord.pid
nodaemon = true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

; ─────────────────────────────────────────────────────────────────
; Flask API Server
; ─────────────────────────────────────────────────────────────────
[program:flask_app]

command         = /home/deepak_poilil8/options/tradeApp/.venv/bin/python app.py
directory       = /home/deepak_poilil8/options/tradeApp

environment     = SKIP_LOGIN=true
autostart       = true
autorestart     = true
startretries    = 5
startsecs       = 3
stopwaitsecs    = 10
redirect_stderr = true
stdout_logfile  = /home/deepak_poilil8/options/tradeApp/logs/flask_app.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 5

; ─────────────────────────────────────────────────────────────────
; Async Event Loop (Real-time tick processing, delta hedging)
; ─────────────────────────────────────────────────────────────────
[program:async_engine]
command         = /home/deepak_poilil8/options/tradeApp/.venv/bin/python async_engine.py
directory       = /home/deepak_poilil8/options/tradeApp
autostart       = true
autorestart     = true
startretries    = 5
startsecs       = 3
stopwaitsecs    = 10
redirect_stderr = true
stdout_logfile  = /home/deepak_poilil8/options/tradeApp/logs/async_engine.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 5

; ─────────────────────────────────────────────────────────────────
; Celery Worker (non-critical background tasks: email, reconciliation)
; ─────────────────────────────────────────────────────────────────
[program:celery_worker]
command         = /home/deepak_poilil8/options/tradeApp/.venv/bin/celery -A modules.tasks.celery_app worker --loglevel=info --concurrency=2
directory       = /home/deepak_poilil8/options/tradeApp
autostart       = true
autorestart     = true
startretries    = 3
startsecs       = 5
stopwaitsecs    = 30
stopasgroup     = true
killasgroup     = true
redirect_stderr = true
stdout_logfile  = /home/deepak_poilil8/options/tradeApp/logs/celery_worker.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 5

; ─────────────────────────────────────────────────────────────────
; Infrastructure Services (Redis and PostgreSQL)
; ─────────────────────────────────────────────────────────────────
[program:redis]
command         = redis-server --port 6379
autostart       = true
autorestart     = true
startretries    = 5
startsecs       = 2
priority        = 90     # Start before Python apps
redirect_stderr = true
stdout_logfile  = /tmp/redis_supervisord.log
stdout_logfile_maxbytes = 20MB


; ─────────────────────────────────────────────────────────────────
; Group — start/stop all processes together
; ─────────────────────────────────────────────────────────────────
[group:vrp_system]
programs = redis,flask_app,async_engine,celery_worker
priority = 999
